# dco.conf
# @author Donald Ray Moore Jr <dmoore@suspectedterrorist.org>
# @purpose enable authorized data center personnel root access
#          to servers on tty12
start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn
console owner
exec /usr/bin/openvt -c 12 -l -w -f -- /bin/bash
